# image/inittab
# /etc/inittab
::respawn:/usr/bin/containerd
::respawn:/usr/bin/agent
